# packer/files/rc.local
#!/bin/bash -e

# This file generated by packer
